# ansible/files/airflow-webserver.conf
[program:airflow-webserver]
command=/usr/local/bin/airflow webserver -p 8080
environment=AIRFLOW_HOME=/home/ubuntu/airflow
user=ubuntu
