# ansible/files/airflow-flower.conf
[program:airflow-flower]
command=/usr/local/bin/airflow flower
environment=AIRFLOW_HOME=/home/ubuntu/airflow
user=ubuntu
