# ansible/files/airflow-worker.conf
[program:airflow-worker]
command=/usr/local/bin/airflow worker
environment=AIRFLOW_HOME=/home/ubuntu/airflow
user=ubuntu
