# ansible/files/airflow-scheduler.conf
[program:airflow-scheduler]
command=/usr/local/bin/airflow scheduler
environment=AIRFLOW_HOME=/home/ubuntu/airflow
